# debian/xrdp.conf
start on runlevel [2345]
stop on runlevel [!2345]

script
	if [ -f /etc/default/xrdp ];then . /etc/default/xrdp;fi
	xrdp -ns
end script
