# debian/xrdp-sesman.conf
start on starting xrdp
stop on stopping xrdp

script
	if [ -f /etc/default/xrdp ];then . /etc/default/xrdp;fi
	xrdp-sesman -ns
end script
